# launchd/com.pandastack.drive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- pandastack-drive on a 4h heartbeat via the audit wrapper (drive-cron.py).
     launchd never runs a 2nd instance of this job while one is live, so overlap is
     impossible without flock. The driver is read-only (Codex sandbox) and proposes
     Linear advances rather than writing them. Install:
       cp launchd/com.pandastack.drive.plist ~/Library/LaunchAgents/
       launchctl load -w ~/Library/LaunchAgents/com.pandastack.drive.plist
     Tune fan-out with PSDRIVE_MAX (default 1 — observe before raising). -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pandastack.drive</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/panda/site/skills/pandastack/scripts/drive-cron.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/panda/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/local/bin</string>
    <key>PSDRIVE_MAX</key>
    <string>1</string>
  </dict>
  <key>StartInterval</key>
  <integer>14400</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/panda/Library/Logs/pandastack-drive/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/panda/Library/Logs/pandastack-drive/launchd.log</string>
</dict>
</plist>
